# supervisor.d/certbot.conf
[program:certbot]
command=certbot --nginx --agree-tos -m %(ENV_EMAIL)s -n -d %(ENV_DOMAINS)s
priority=206
autostart = true
startsecs = 0
startretries = 1
autorestart = false
stdout_logfile= /dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
username=nginx
